# deploy/kalshi-hedger.service
[Unit]
Description=Kalshi Hedger Pool Bot
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=kalshi
Group=kalshi
WorkingDirectory=/opt/kalshi-trading-bot
ExecStart=/opt/kalshi-trading-bot/venv/bin/python /opt/kalshi-trading-bot/hedger.py

# ---- Environment ----
EnvironmentFile=-/opt/kalshi-trading-bot/.env

# ---- Restart policy ----
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# ---- Security hardening ----
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/kalshi-trading-bot
PrivateTmp=true

# ---- Logging ----
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kalshi-hedger

[Install]
WantedBy=multi-user.target
